# launchd/io.one2b.brief-watchman.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.one2b.brief-watchman</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jasonpeterstevens/.gmail-mcp-intel/brief-watchman.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jasonpeterstevens/.gmail-mcp-intel</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/brief-watchman-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/brief-watchman-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jasonpeterstevens/.npm-global/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/jasonpeterstevens</string>
    </dict>
</dict>
</plist>
